# tailscale/rootfs/etc/s6-overlay/s6-rc.d/protect-subnets/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Prevent local subnets to be routed toward the tailnet
# ==============================================================================

protect-subnet-routes

# runs scripts in /etc/NetworkManager/dispatcher.d
# --debug is used to prevent logging to syslog (HA cli)
exec /usr/libexec/nm-dispatcher --persist --debug > /dev/null
