# tailscale/rootfs/etc/s6-overlay/s6-rc.d/share-homeassistant/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Enables Tailscale Serve or Funnel feature to share Home Assistant
# ==============================================================================

readonly WAIT_DELAY=5   # 5s
readonly WAIT_COUNT=18  # 18*5s = 90s = 1.5m

declare wait_counter=0
declare curl_result

# Validate share_homeassistant value
if ! bashio::config.equals 'share_homeassistant' 'serve' && \
  ! bashio::config.equals 'share_homeassistant' 'funnel'
then
  bashio::exit.nok "Invalid value '$(bashio::config 'share_homeassistant')' for share_homeassistant. Must be either 'serve' or 'funnel'"
fi

# Check if Tailscale HTTPS is enabled
if ! /opt/tailscale status --self=true --peers=false --json \
  | jq -rce '.Self.CapMap | has("https")' > /dev/null;
then
  bashio::exit.nok "Tailscale's HTTPS support is disabled"
fi

# Checking if SSL is used
if bashio::var.true "$(bashio::core.ssl)"; then
  bashio::exit.nok "Tailscale's HTTPS support is enabled, but Home Assistant is not accessible through plain HTTP connection"
fi

# Wait a bit for HA to be available during startup
while (( 200 != (curl_result=$(curl -s -o /dev/null -w "%{http_code}" "http://127.0.0.1:$(bashio::core.port)")) )); do
  if (( wait_counter++ == WAIT_COUNT )); then
    break
  fi
  bashio::log.info "Waiting for Home Assistant to be ready..."
  sleep $WAIT_DELAY
done
if (( wait_counter != 0 && curl_result == 200)); then
  bashio::log.info "Home Assistant is ready"
fi

if (( 200 != curl_result )); then
  # Warn that we can't test Home Assistant's HTTP reverse proxy configuration
  # We emit only a warning to let the add-on start, maybe this is the only connection to access the device, better to start than not
  # Though starting tailscale serve without HA won't fail, but will fill the logs with messages forever in each ~10s
  bashio::log.warning "Home Assistant is not accessible currently, unable to test the connection to Home Assistant as reverse proxy"
else
  # Test Home Assistant's HTTP reverse proxy configuration
  if (( 200 != $(curl -s -o /dev/null -w "%{http_code}" "http://127.0.0.1:$(bashio::core.port)" -H "X-Forwarded-For: 127.0.0.1") )); then
    bashio::exit.nok \
      "Unable to connect to Home Assistant as reverse proxy." \
      "Please check your configuration based on the add-on's documentation under \"Option: share_homeassistant\"."
  fi
fi

# Check if Funnel is available
if bashio::config.equals 'share_homeassistant' 'funnel'; then
  if ! /opt/tailscale status --self=true --peers=false --json \
    | jq -rce '.Self.CapMap | has("funnel")' > /dev/null;
  then
    bashio::exit.nok "Tailscale's Funnel support is disabled"
  fi
fi

# This service can wait for HA for minutes, let notify S6 when we are really starting
echo "" >&3

# Set up serve or funnel
exec /opt/tailscale $(bashio::config 'share_homeassistant') --https=$(bashio::config 'share_on_port' '443') --set-path=/ "http://127.0.0.1:$(bashio::core.port)"
