# tailscale/rootfs/etc/s6-overlay/s6-rc.d/taildrop/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Fetches files send via Taildrop
# ==============================================================================

# Ensure the directory exists
mkdir -p /share/taildrop

# Fetch files
exec /opt/tailscale file get --verbose --loop "/share/taildrop"
