# tailscale/rootfs/etc/s6-overlay/s6-rc.d/init-magicdns-proxies/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
export LOG_FD
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Remove forwarding
# ==============================================================================

magicdns-ingress-proxy-forwarding remove drop
