# tailscale/rootfs/etc/s6-overlay/s6-rc.d/init-magicdns-proxies/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Runs the dnsmasq proxies initialization
# ==============================================================================

readonly DEFAULT_LOGIN_SERVER="controlplane.tailscale.com"
readonly LETSENCRYPT_API="acme-v02.api.letsencrypt.org"
readonly DNSMASQ_BLACK_WHITE_LIST_LOCATION="/etc/dnsmasq-black-white-list"

declare login_server="${DEFAULT_LOGIN_SERVER}"
declare -a black_white_list=()

# We have to be able to determine login_server from this address
if ! login_server=$(awk -F[/:] '{print $4}' <<<$(bashio::config "login_server")) || \
  ! bashio::var.has_value "${login_server}"
then
  bashio::log.error "Determining host name from '$(bashio::config "login_server")' has failed"
  bashio::exit.nok
fi
black_white_list+=(${login_server})

# If serve or funnel is used, resolve letsencrypt's api also
if bashio::config.equals 'proxy' 'true'; then
  black_white_list+=(${LETSENCRYPT_API})
fi

printf "%s" "${black_white_list[@]/%/$'\n'}" > ${DNSMASQ_BLACK_WHITE_LIST_LOCATION}

# This is necessary to prevent accessing MagicDNS before the ingress proxy starts up
# The ingress proxy will remove these entries on startup
magicdns-ingress-proxy-forwarding setup drop
