# tailscale/rootfs/etc/s6-overlay/s6-rc.d/web/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Runs tailscale web interface
# ==============================================================================
declare -a options

bashio::log.info 'Starting Tailscale web...'

# Get random port
options+=(--listen 127.0.0.1:25899)

# Use readonly webui mode
# Note: can be removed later, after most of the add-on config options are available through the webui
options+=(--readonly)

# Run Tailscale
exec /opt/tailscale web "${options[@]}"
