# tailscale/rootfs/etc/s6-overlay/s6-rc.d/magicdns-egress-proxy/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Runs the MagicDNS egress proxy
# ==============================================================================

readonly DNSMASQ_EGRESS_ADDRESS_IPV4="127.100.100.100"
readonly DNSMASQ_EGRESS_PORT=53
readonly DNSMASQ_BLACK_WHITE_LIST_LOCATION="/etc/dnsmasq-black-white-list"

declare hassio_dns_ipv4
declare -a white_list

bashio::log.info "Starting MagicDNS egress proxy..."

function dig_hassio_dns() {
  local type="$1"
  dig dns.local.hass.io "${type}" +short | { grep -Ev '^;|\.$|^$' || true ;}
}

# Hassio DNS's IP addresses
if ! hassio_dns_ipv4=$(dig_hassio_dns A) || \
  bashio::var.is_empty "${hassio_dns_ipv4}"
then
  bashio::log.error "Failed to resolve Home Assistant's DNS address"
  bashio::exit.nok
fi

# White-list
readarray -t white_list < ${DNSMASQ_BLACK_WHITE_LIST_LOCATION}

# We need to delay the starting of the dependent services until the conf file is written
echo "nameserver ${DNSMASQ_EGRESS_ADDRESS_IPV4}" > /etc/resolv.dnsmasq.conf
echo "" >&3

# This DNS forwards the white_list to HA's DNS, otherwise replies REFUSED for everything
# It must run on port 53 to be able to specify it in a resolv.conf
exec dnsmasq --no-hosts --no-resolv --keep-in-foreground --log-facility='-' --cache-size=0 \
  --listen-address=${DNSMASQ_EGRESS_ADDRESS_IPV4} --bind-interfaces \
  --port=${DNSMASQ_EGRESS_PORT} \
  $(printf -- " --server=/%s/${hassio_dns_ipv4}" "${white_list[@]}")
