# tailscale/rootfs/etc/s6-overlay/s6-rc.d/init-magicdns-proxies-upstream-list/run
#!/command/with-contenv bashio
# shellcheck shell=bash
export LOG_FD
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Runs the dnsmasq proxies upstream black/white list initialization
# ==============================================================================

# Note: This service runs when userspace_networking is disabled, and
# - when accept_dns is enabled, adds login server, log server and letsencrypt api to the black_white_list,
#   to prevent DNS query loops (both egress and ingress proxies are running)
# - when accept_dns is disabled, adds tailnet to the black_white_list,
#   to prevent logging SERVFAIL by 100.100.100.100 for non-tailnet queries (only ingress proxy is running)

readonly MAGIC_DNS_IPV4="100.100.100.100"
readonly MAGIC_DNS_IPV6="fd7a:115c:a1e0::53"

readonly DEFAULT_LOGIN_SERVER="controlplane.tailscale.com"
readonly LOG_SERVER="log.tailscale.com"
readonly LETSENCRYPT_API="acme-v02.api.letsencrypt.org"
readonly DNSMASQ_BLACK_WHITE_LIST_LOCATION="/etc/dnsmasq-black-white-list"

declare login_server="${DEFAULT_LOGIN_SERVER}"
declare -a black_white_list=()
declare tailnet

if bashio::config.true "accept_dns"; then
  # We have to be able to determine login_server from this address
  if ! login_server=$(awk -F[/:] '{print $4}' <<<"$(bashio::config "login_server")") || \
    ! bashio::var.has_value "${login_server}"
  then
    bashio::exit.nok "Determining host name from '$(bashio::config "login_server")' has failed"
  fi
  black_white_list+=("${login_server}")

  # When log upload is enabled, resolve log server also
  if bashio::debug; then
    black_white_list+=("${LOG_SERVER}")
  fi

  # If serve or funnel is used, resolve letsencrypt's api also
  if ! bashio::config.equals 'share_homeassistant' 'disabled'; then
    black_white_list+=("${LETSENCRYPT_API}")
  fi
else
  # Add the tailnet domain as white list
  if ! tailnet=$( \
      /opt/tailscale status --json --peers=false --self=false \
      | jq -re '.MagicDNSSuffix') || \
    ! bashio::var.has_value "${tailnet}"
  then
    bashio::exit.nok "Determining tailnet name has failed"
  fi
  black_white_list+=("${tailnet}")
fi

printf "%s" "${black_white_list[@]/%/$'\n'}" > "${DNSMASQ_BLACK_WHITE_LIST_LOCATION}"
