# tailscale/rootfs/etc/s6-overlay/s6-rc.d/funnel/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Enables Tailscale Funnel feature
# ==============================================================================

# Check if Tailscale HTTPS is enabled
if ! /opt/tailscale status --self=true --peers=false --json \
  | jq -rce '.Self.CapMap | has("https")' > /dev/null;
then
  bashio::log.error "Tailscale's HTTPS support is disabled"
  bashio::exit.nok
fi

# Check if Funnel is available
if ! /opt/tailscale status --self=true --peers=false --json \
  | jq -rce '.Self.CapMap | has("funnel")' > /dev/null;
then
  bashio::log.error "Tailscale's Funnel support is disabled"
  bashio::exit.nok
fi

# Set up funnel
exec /opt/tailscale funnel --https=443 --set-path=/ "http://127.0.0.1:$(bashio::core.port)"
