# tailscale/rootfs/etc/s6-overlay/s6-rc.d/mss-clamping/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Clamp the MSS to the MTU
# ==============================================================================

readonly TAILSCALE_INTERFACE="tailscale0"
readonly CLAMPING_IPTABLES_OPTIONS="-p tcp -m tcp --tcp-flags SYN,RST SYN -j TCPMSS --clamp-mss-to-pmtu"

# In case of non userspace networking, clamp the MSS to the MTU for tailscale's interface
# If user later enables subnet routing for site-to-site networking, this config is already there
# Source: https://tailscale.com/kb/1214/site-to-site/
function setup_clamping() {
  local cmd="$1"
  local ip_version="$2"

  bashio::log.info "Clamping the MSS to the MTU for interface ${TAILSCALE_INTERFACE}," \
    "to support site-to-site networking better (${ip_version})"
  if ${cmd} -t mangle -S FORWARD \
    | grep -Eq "^-A FORWARD -o tailscale\d ${CLAMPING_IPTABLES_OPTIONS}$"
  then
    bashio::log.notice "Clamping is already set (${ip_version})"
  else
    if ! ${cmd} -t mangle -A FORWARD -o ${TAILSCALE_INTERFACE} ${CLAMPING_IPTABLES_OPTIONS}; then
      bashio::log.warning "Setting up clamping is unsuccessful (${ip_version})"
    fi
  fi
}

setup_clamping "iptables" "IPv4"
setup_clamping "ip6tables" "IPv6"
