# tailscale/rootfs/etc/services.d/cert/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Runs tailscale cert periodically
# ==============================================================================

# Request TLS certificate if certificate_tailnet_name is set
tailscale_cert() {
  if bashio::config.has_value 'certificate_tailnet_name'; then
    if ! bashio::fs.directory_exists "/ssl/tailscale"; then
        mkdir -p "/ssl/tailscale"
    fi
    /opt/tailscale cert \
      --cert-file "/ssl/tailscale/fullchain.pem" \
      --key-file "/ssl/tailscale/privkey.pem" \
      "$(bashio::info.hostname).$(bashio::config 'certificate_tailnet_name')"
  fi
}

cert_err='Error running "tailscale cert". Visit https://tailscale.com/kb/1153/enabling-https/ to setup HTTPS.'

if ! tailscale_cert; then
  bashio::exit.nok "$cert_err"
fi

bashio::log.info 'Starting Tailscale TLS certificate fetcher...'

while true;
do
  # Fetch certs once a day
  sleep 86400
  # Fetch only if expiration is within 1 week
  if ! openssl x509 -checkend 604800 -noout -in /ssl/tailscale/fullchain.pem > /dev/null; then
    if ! tailscale_cert; then
      bashio::log.error "$cert_err"
    fi
  fi
done
