# tailscale/rootfs/etc/s6-overlay/s6-rc.d/tailscaled/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Take down the S6 supervision tree when Tailscale fails
# ==============================================================================
readonly exit_code_container=$(</run/s6-linux-init-container-results/exitcode)
readonly exit_code_service="${1}"
readonly exit_code_signal="${2}"
readonly service="tailscaled"

declare tailscaled_gid

# Remove DSCP setting
for tailscaled_gid in $( \
  iptables -t mangle -S OUTPUT \
    | { grep -E '^-A OUTPUT -m owner --gid-owner \d* -j DSCP --set-dscp' || true ;} \
    | sed -nr 's/^.*?--gid-owner\s(\d+)\s.*$/\1/p')
do
  bashio::log.info "Removing DSCP setting for tailscaled (IPv4)"
  iptables -t mangle -D OUTPUT -m owner --gid-owner $tailscaled_gid -j DSCP --set-dscp $(bashio::config "dscp")
done
for tailscaled_gid in $( \
  ip6tables -t mangle -S OUTPUT \
    | { grep -E '^-A OUTPUT -m owner --gid-owner \d* -j DSCP --set-dscp' || true ;} \
    | sed -nr 's/^.*?--gid-owner\s(\d+)\s.*$/\1/p')
do
  bashio::log.info "Removing DSCP setting for tailscaled (IPv6)"
  ip6tables -t mangle -D OUTPUT -m owner --gid-owner $tailscaled_gid -j DSCP --set-dscp $(bashio::config "dscp")
done

bashio::log.info \
  "Service ${service} exited with code ${exit_code_service}" \
  "(by signal ${exit_code_signal})"

if [[ "${exit_code_service}" -eq 256 ]]; then
  if [[ "${exit_code_container}" -eq 0 ]]; then
    echo $((128 + $exit_code_signal)) > /run/s6-linux-init-container-results/exitcode
  fi
  [[ "${exit_code_signal}" -eq 15 ]] && exec /run/s6/basedir/bin/halt
elif [[ "${exit_code_service}" -ne 0 ]]; then
  if [[ "${exit_code_container}" -eq 0 ]]; then
    echo "${exit_code_service}" > /run/s6-linux-init-container-results/exitcode
  fi
  exec /run/s6/basedir/bin/halt
fi
