# tailscale/rootfs/etc/s6-overlay/s6-rc.d/tailscaled/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Runs tailscale
# ==============================================================================
declare -a options

bashio::log.info 'Starting Tailscale...'

options+=(--state=/data/tailscaled.state)
options+=(--tun=userspace-networking)

# Run Tailscale
exec /opt/tailscaled "${options[@]}"
