# tailscale/rootfs/etc/s6-overlay/s6-rc.d/proxy/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Enables Tailscale Proxy feature
# ==============================================================================

declare domain

# Check if Tailscale HTTPS is enabled
if ! /opt/tailscale status --self=true --peers=false --json \
  | jq -rce '.Self.CapMap | has("https")' > /dev/null;
then
  bashio::log.error "Tailscale's HTTPS support is disabled. Therefore, can't start the add-on's Tailscale Proxy functionality."
  bashio::exit.nok
fi

domain=$(/opt/tailscale status --self=true --peers=false --json | jq -rc ".CertDomains[0]")

# Checking if SSL is used
if bashio::var.true "$(bashio::core.ssl)"; then
  bashio::log.error "Tailscale's HTTPS support is enabled, but Home Assistant is not accessible through plain HTTP connection"
  bashio::exit.nok
fi

# Set up proxy
if ! /opt/tailscale serve https:443 / "http://127.0.0.1:$(bashio::core.port)"; then
  bashio::log.error "Unable to configure Tailscale Proxy"
  bashio::exit.nok
fi
bashio::log.info "Tailscale Proxy is enabled:"
bashio::log.info "  Your Home Assistant instance is available within your Tailnet VPN at"
bashio::log.info "  https://${domain}"
