# tailscale/rootfs/etc/s6-overlay/s6-rc.d/proxy/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Enables Tailscale Proxy feature
# ==============================================================================

declare domain

# Tailscale serve and funnel configuration is persistent between restarts
# Disable previously configured tailscale serve functionality with deleting the tailscale serve configuration
if ! /opt/tailscale serve reset; then
  bashio::log.error "Unable to reset Tailscale Proxy settings"
  bashio::exit.nok
fi

# Enable proxy service only when it has been explicitly enabled
if bashio::config.true 'proxy'; then

  # Check if Tailscale HTTPS is enabled
  if ! /opt/tailscale status --self=true --peers=false --json \
    | jq -rce '.Self.CapMap | has("https")' > /dev/null;
  then
    bashio::log.error "Tailscale's HTTPS support is disabled"
    bashio::exit.nok
  fi

  domain=$(/opt/tailscale status --self=true --peers=false --json | jq -rc ".CertDomains[0]")

  # Checking if SSL is used
  if bashio::var.true "$(bashio::core.ssl)"; then
    bashio::log.error "Tailscale's HTTPS support is enabled, but Home Assistant is not accessible through plain HTTP connection"
    bashio::exit.nok
  fi

  if (( 200 != $(curl -s -o /dev/null -w "%{http_code}" "http://127.0.0.1:$(bashio::core.port)") )); then
    # Warn that we can't test Home Assistant's HTTP reverse proxy configuration
    bashio::log.warning "Home Assistant is not accessible, unable to test the connection to Home Assistant as reverse proxy"
  else
    # Test Home Assistant's HTTP reverse proxy configuration
    if (( 200 != $(curl -s -o /dev/null -w "%{http_code}" "http://127.0.0.1:$(bashio::core.port)" -H "X-Forwarded-For: 127.0.0.1") )); then
      bashio::log.error "Unable to connect to Home Assistant as reverse proxy"
      bashio::log.error "Please check your configuration based on the add-on's documentation under \"Option: proxy\""
      bashio::exit.nok
    fi
  fi

  # Set up proxy
  if ! /opt/tailscale serve https:443 / "http://127.0.0.1:$(bashio::core.port)"; then
    bashio::log.error "Unable to configure Tailscale Proxy"
    bashio::exit.nok
  fi
  bashio::log.info "Tailscale Proxy is enabled:"
  bashio::log.info "  Your Home Assistant instance is available within your tailnet VPN at"
  bashio::log.info "  https://${domain}"

fi
