# tailscale/rootfs/etc/s6-overlay/s6-rc.d/local-network/up
/etc/s6-overlay/s6-rc.d/local-network/run
