# tailscale/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Configures NGINX for use with the Tailscale web
# ==============================================================================


# Funnel protection is disabled by default
funnel_protection=false
if bashio::config.true "funnel_protection";
then
  funnel_protection=true
fi

# Generate Ingress configuration
bashio::var.json \
    interface "$(bashio::addon.ip_address)" \
    port "^$(bashio::addon.ingress_port)" \
    | tempio \
        -template /etc/nginx/templates/ingress.gtpl \
        -out /etc/nginx/servers/ingress.conf

bashio::var.json \
    port "$(bashio::core.port)" \
    funnel_protection "^${funnel_protection}" \
    | tempio \
        -template /etc/nginx/templates/homeassistant.gtpl \
        -out /etc/nginx/servers/homeassistant.conf
