# tailscale/rootfs/etc/services.d/tailscaled/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Runs tailscale
# ==============================================================================
declare -a options

bashio::log.info 'Starting Tailscale...'

options+=(--state=/data/tailscaled.state)
options+=(--tun=userspace-networking)

# Trigger post process to run, after stuff got connected
./post &

# Run Tailscale
exec /opt/tailscaled "${options[@]}"
