# tailscale/rootfs/etc/s6-overlay/s6-rc.d/post-tailscaled/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Runs after the machine has been logged in into the Tailscale network
# ==============================================================================
declare -a addresses=()
declare -a routes=()
declare ipinfo
declare tags

function appendarray() {
  local -n array=${1}
  readarray -t -O "${#array[@]}" array
}

# Find interfaces and matching addresses from which we can extract routes to be advertised
for interface in $(bashio::network.interfaces); do
  appendarray addresses < <(bashio::network.ipv4_address "${interface}")
  appendarray addresses < <(bashio::network.ipv6_address "${interface}")
done

# Extract routes to be advertised
for address in "${addresses[@]}"; do
  if bashio::var.has_value "${address}"; then
    # Skip local link addresses
    if [[ "${address:0:6}" == "fe80::" ]] || [[ "${address:0:8}" == "169.254." ]];
    then
      continue
    fi

    # Skip if forwarding for the address family is disabled
    if [[ "${address}" =~ .*:.* ]];
    then
      [[ $(</proc/sys/net/ipv6/conf/all/forwarding) -eq 0 ]] && continue
    else
      [[ $(</proc/sys/net/ipv4/ip_forward) -eq 0 ]] && continue
    fi

    ipinfo="$(/usr/bin/ipcalc --json "${address}")"
    routes+=("$(bashio::jq "${ipinfo}" '.NETWORK + "/" + .PREFIX')")
  fi
done
# Remove duplicate entries
readarray -t routes < <(printf "%s\n" "${routes[@]}" | sort -u)

# Get configured tags
tags=$(bashio::config "tags//[] | join(\",\")" "")

# Get server to login against
login_server=$(bashio::config 'login_server' 'https://controlplane.tailscale.com')

# Wait for socket to be available
while ! bashio::fs.socket_exists "/var/run/tailscale/tailscaled.sock";
do
  sleep 1
done

# Wait for the network to be available and logged in
while true;
do
  if /opt/tailscale status --json --peers=false --self=false \
    | jq --exit-status '.BackendState == "Running" or .BackendState == "NeedsLogin"' > /dev/null;
  then
    IFS=","
    /opt/tailscale up \
      --hostname "$(bashio::info.hostname)" \
      --advertise-exit-node \
      --accept-routes \
      --advertise-routes="${routes[*]}" \
      --advertise-tags="${tags}" \
      --login-server="${login_server}"

    bashio::exit.ok
  fi
  # Well... wait a bit more
  sleep 2
done
