# tailscale/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Runs the Nginx daemon
# ==============================================================================

# Wait for the Tailscale webinterface to become available
bashio::net.wait_for 25899 localhost 900

bashio::log.info "Starting NGinx..."
exec nginx
