# tailscale/rootfs/etc/s6-overlay/s6-rc.d/taildrive/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Exposes Home Assistant directories over Taildrive
# ==============================================================================

# Read configured directories from the add-on configuration.
if ! mapfile -t configured_dirs < <(bashio::config "taildrive"); then
    bashio::exit.nok "Error reading configured directories from add-on configuration."
fi

# Read currently shared directories into an array from tailscale drive list.
# The output of which looks like:
# name    path    as
# ----    ----    --
# config  /config  config
if ! mapfile -t shared_dirs < <(/opt/tailscale drive list | tail -n +3 | awk '{print $1}'); then
    bashio::exit.nok "Error reading shared directories from Tailscale."
fi

# If a directory is configured but not shared, share it.
for dir in "${configured_dirs[@]}"; do
    if ! printf '%s\0' "${shared_dirs[@]}" | grep -Fxqz -- "${dir}"; then
        # Share the directory
        /opt/tailscale drive share "${dir}" "/${dir}"
    fi
done

# If a directory is shared but not configured, unshare it.
for dir in "${shared_dirs[@]}"; do
    if ! printf '%s\0' "${configured_dirs[@]}" | grep -Fxqz -- "${dir}"; then
        # Unshare the directory
        /opt/taildrive drive unshare "${dir}"
    fi
done