# tailscale/rootfs/etc/s6-overlay/s6-rc.d/taildrive/run
#!/command/with-contenv bashio
# shellcheck shell=bash
export LOG_FD
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Exposes Home Assistant directories over Taildrive
# ==============================================================================

# These are not arrays, but newline separated lists (single strings)
declare all_share_names
declare configured_share_names
declare active_share_names

declare share_name

# Read all possible and configured directories from the add-on configuration.
if ! all_share_names=$(bashio::config 'taildrive' | jq -rc 'keys_unsorted[]') || \
  ! configured_share_names=$(bashio::config 'taildrive' | jq -rc 'to_entries[] | select(.value) | .key')
then
  bashio::exit.nok "Error reading configured Taildrive directories from add-on configuration."
fi

# Check if Taildrive is available
if ! /opt/tailscale status --self=true --peers=false --json \
  | jq -rce '.Self.CapMap | has("drive:share")' > /dev/null;
then
  # If Taildrive is not available but any share is configured, that is an error
  if (( 0 < $(wc -w <<< "${configured_share_names}") )); then
    bashio::exit.nok "Tailscale's Taildrive is disabled but directories are configured for sharing"
  fi
else
  # Read currently shared directories into an array from tailscale drive list.
  # The output of which looks like:
  # name      path       as
  # ------    -------    ----
  # config    /config    root
  if ! active_share_names=$(/opt/tailscale drive list | tail -n +3 | awk '{print $1}'); then
    bashio::exit.nok "Error reading shared Taildrive directories from Tailscale."
  fi

  ln -s "/homeassistant" "/config"

  # If a directory is configured but not shared, share it.
  # Note: do not quote configured_share_names in the for loop,
  # as we want to iterate over words (lines), this is not an array
  for share_name in ${configured_share_names}; do
    if ! grep -Fxq -- "${share_name}" <<< "${active_share_names}"; then
      bashio::log.info "Taildrive: sharing /${share_name}"
      if ! /opt/tailscale drive share "${share_name}" "/${share_name}"; then
        bashio::exit.nok "Failed to share /${share_name}"
      fi
    else
      bashio::log.info "Taildrive: /${share_name} is already shared"
    fi
  done

  # If a directory is shared but not configured, and managed by this add-on, unshare it.
  # Note: do not quote active_share_names in the for loop,
  # as we want to iterate over words (lines), this is not an array
  for share_name in ${active_share_names}; do
    if ! grep -Fxq -- "${share_name}" <<< "${configured_share_names}"; then
      if grep -Fxq -- "${share_name}" <<< "${all_share_names}"; then
        bashio::log.info "Taildrive: unsharing /${share_name}"
        if ! /opt/tailscale drive unshare "${share_name}"; then
          bashio::exit.nok "Failed to unshare /${share_name}"
        fi
      else
        bashio::log.info "Taildrive: skipping to unshare '${share_name}'"
      fi
    fi
  done
fi
