# tailscale/rootfs/etc/s6-overlay/s6-rc.d/dnsmasq/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Runs the dnsmasq daemon
# ==============================================================================

bashio::log.info "Starting dnsmasq..."
# This is a dummy DNS to suppress tailscaled warnings about not configured upstream on each DNS query
# It answers REFUSED for everything
touch /etc/dnsmasq.conf
echo "nameserver 127.0.0.1:25798" > /etc/resolv.dnsmasq.conf
exec dnsmasq --no-hosts --no-resolv --keep-in-foreground --log-facility='-' --listen-address=127.0.0.1 --port=25798 --local-service=host
