# tailscale/rootfs/etc/s6-overlay/s6-rc.d/dnsmasq/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Runs the dnsmasq daemon
# ==============================================================================

readonly DNSMASQ_ADDRESS=127.52.52.52

bashio::log.info "Starting dnsmasq..."
# This is a dummy DNS to suppress tailscaled warnings about not configured upstream on each DNS query
# It answers REFUSED for everything
# It must run on port 53 to be able to specify it in a resolv.conf
echo "nameserver ${DNSMASQ_ADDRESS}" > /etc/resolv.dnsmasq.conf
exec dnsmasq --no-hosts --no-resolv --keep-in-foreground --log-facility='-' --listen-address=${DNSMASQ_ADDRESS} --port=53 --local-service=host
