# tailscale/rootfs/etc/s6-overlay/s6-rc.d/magicdns-ingress-proxy/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Runs the MagicDNS ingress proxy
# ==============================================================================

readonly MAGIC_DNS_IPV4="100.100.100.100"

readonly DNSMASQ_INGRESS_PORT=53
readonly DNSMASQ_BLACK_WHITE_LIST_LOCATION="/etc/dnsmasq-black-white-list"

declare tailscale_address_ipv4
declare tailscale_address_ipv6
declare -a listen_addresses
declare -a black_list

bashio::log.info "Starting MagicDNS ingress proxy..."

# Tailscale's local IP addresses
if ! tailscale_address_ipv4=$(/opt/tailscale ip -4) || \
  ! tailscale_address_ipv6=$(/opt/tailscale ip -6) || \
  bashio::var.is_empty "${tailscale_address_ipv4}" && bashio::var.is_empty "${tailscale_address_ipv6}"
then
  bashio::log.error "Failed to retrieve Tailscale's local address"
  bashio::exit.nok
fi

# Listen addresses
if bashio::var.has_value "${tailscale_address_ipv4}"; then
  listen_addresses+=(${tailscale_address_ipv4})
fi
if bashio::var.has_value "${tailscale_address_ipv6}"; then
  listen_addresses+=(${tailscale_address_ipv6})
fi

# Black-list
readarray -t black_list < ${DNSMASQ_BLACK_WHITE_LIST_LOCATION}

magicdns-ingress-proxy-forwarding setup forwarding
magicdns-ingress-proxy-forwarding remove drop

# We need to delay the starting of the dependent services until iptables are configured
echo "" >&3

# This DNS replies NXDOMAIN for the black_list, otherwise forwards queries to MagicDNS
exec dnsmasq --no-hosts --no-resolv --keep-in-foreground --log-facility='-' --cache-size=0 \
  $(printf -- " --listen-address=%s" "${listen_addresses[@]}") --bind-dynamic \
  --port=${DNSMASQ_INGRESS_PORT} \
  --server=${MAGIC_DNS_IPV4} $(printf -- " --server=/%s/" "${black_list[@]}")
