# tailscale/rootfs/etc/s6-overlay/s6-rc.d/magicdns-ingress-proxy/run
#!/command/with-contenv bashio
# shellcheck shell=bash
export LOG_FD
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Runs the MagicDNS ingress proxy
# ==============================================================================

# Note: This service runs when userspace_networking is disabled, and
# - when accept_dns is enabled, it treats black_white_list as black_list, and
#   returns NXDOMAIN for anything that is on the list
# - when accept_dns is disabled, it treats black_white_list as white_list, and
#   returns NXDOMAIN for anything that is NOT on the list

source /usr/lib/trace.sh

readonly MAGIC_DNS_IPV4="100.100.100.100"

readonly DNSMASQ_INGRESS_PORT=53
readonly DNSMASQ_BLACK_WHITE_LIST_LOCATION="/etc/dnsmasq-black-white-list"

declare tailscale_address_ipv4
declare tailscale_address_ipv6
declare -a black_list
declare -a white_list
declare domain
declare -a options

bashio::log.info "Starting MagicDNS ingress proxy..."

options+=(--no-hosts)
options+=(--no-resolv)
options+=(--conf-file=/dev/null)
options+=(--keep-in-foreground)
options+=(--log-facility='-')
options+=(--cache-size=0)

# Tailscale's local IP addresses
if ! tailscale_address_ipv4=$(/opt/tailscale ip -4) || \
  ! tailscale_address_ipv6=$(/opt/tailscale ip -6) || \
  bashio::var.is_empty "${tailscale_address_ipv4-}" && bashio::var.is_empty "${tailscale_address_ipv6-}"
then
  bashio::exit.nok "Failed to retrieve Tailscale's local address"
fi

# Listen addresses
if bashio::var.has_value "${tailscale_address_ipv4-}"; then
  options+=(--listen-address=${tailscale_address_ipv4})
fi
if bashio::var.has_value "${tailscale_address_ipv6-}"; then
  options+=(--listen-address=${tailscale_address_ipv6})
fi

options+=(--bind-dynamic)
options+=(--port=${DNSMASQ_INGRESS_PORT})

if bashio::config.true "accept_dns"; then
  # Black-list
  readarray -t black_list < "${DNSMASQ_BLACK_WHITE_LIST_LOCATION}"

  # Forward everything to MagicDNS, except the black_list
  options+=(--server=${MAGIC_DNS_IPV4})
  for domain in "${black_list[@]}"; do
    options+=(--server=/${domain}/)
  done
else
  # White-list
  readarray -t white_list < "${DNSMASQ_BLACK_WHITE_LIST_LOCATION}"

  # Return NXDOMAIN for everything, except the white_list
  options+=(--address=/#/)
  for domain in "${white_list[@]}"; do
    options+=(--server=/${domain}/${MAGIC_DNS_IPV4})
  done
fi

if bashio_custom::trace; then
  options+=(--log-queries)
  options+=(--log-debug)
fi

magicdns-ingress-proxy-forwarding setup forwarding
magicdns-ingress-proxy-forwarding remove drop

# We need to delay the starting of the dependent services until iptables are configured
echo "" >&3

# This DNS replies NXDOMAIN for the black_list, otherwise forwards everything to MagicDNS
exec dnsmasq "${options[@]}"
